# modules/install-cfssl/cfssl.path
[Path]
PathExists=/etc/sysconfig/cfssl.conf

[Install]
WantedBy=multi-user.target
